# udev/81-input-a2dp-autoconnect.rules
SUBSYSTEM=="input", RUN+="/bin/sh -c '~/pipod/a2dp-autoconnect'" OPTIONS="last_rule"
